# rootfs/etc/s6-overlay/s6-rc.d/svc-php/run
#!/usr/bin/with-contenv sh
# shellcheck shell=sh
exec s6-setuidgid rtorrent php-fpm83 -F